# autodeploy/gsloth.conf
# /etc/init/gsloth.conf

description "Autodeployed Gearsloth"
author      "Antti <antti@meetin.gs>"

start on    (local-filesystems and net-device-up IFACE!=lo)
stop  on    shutdown

respawn
respawn limit 5 5

setuid nobody
setgid nogroup

script
  set -e
  FIFO=`mktemp --dry-run`
  mkfifo $FIFO
  (logger -t sloth -p local7.info <$FIFO &)
  exec 2>$FIFO
  exec >$FIFO
  rm $FIFO

  . /etc/autodeploy.conf

  exec /usr/local/bin/gearslothd
end script
